# systemd/www-print-https.socket
[Socket]
ListenStream=443
BindIPv6Only=both
FileDescriptorName=https
Service=www-print.service

[Install]
WantedBy=sockets.target
