# systemd/www-print-http.socket
[Socket]
ListenStream=80
BindIPv6Only=both
FileDescriptorName=http
Service=www-print.service

[Install]
WantedBy=sockets.target
